# scheduler/crontab
# Script used on Linux to start pipeline

* * * * * /mnt/e/python-dev-test/scheduler/census_bureau_pipeline.sh